# systemd/user/dunst.service
[Unit]
Description=Launch Dunst notification service

[Service]
ExecStart=/usr/bin/dunst

[Install]
WantedBy=graphical.target
