# systemd/user/flameshot.service
[Unit]
Description=Launch Flameshot

[Service]
ExecStart=/usr/bin/flameshot

[Install]
WantedBy=graphical.target
