# systemd/user/mbsync.service
[Unit]
Description=Sync IMAP folders

[Service]
ExecStart=/usr/bin/mbsync -a

[Install]
WantedBy=graphical.target
