# systemd/user/mbsync@.timer
[Unit]
Description=Mbsync timer

[Timer]
OnBootSec=0m
OnUnitActiveSec=5m
Unit=mbsync@%i.service

[Install]
WantedBy=timers.target
