# systemd/user/mbsync@.service
[Unit]
Description=Mbsync service

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync %i

[Install]
WantedBy=default.target
