# systemd/user/polybar.service
[Unit]
Description=(Re)Launch Polybar

[Service]
Type=forking
ExecStart=%h/.config/polybar/launch.sh

[Install]
WantedBy=graphical.target
