# systemd/user/xmodmap.service
[Unit]
Description=Launch XModmap service

[Service]
ExecStart=/usr/bin/xmodmap %h/.Xmodmap

[Install]
WantedBy=graphical.target
