# systemd/user/xss-lock.service
[Unit]
Description=Launch XSS Lock service

[Service]
ExecStart=/usr/bin/xss-lock -- %h/.local/bin/lock

[Install]
WantedBy=graphical.target
